# .config/systemd/user/backup-packages.service
[Unit]
Description=Backup list of installed packages

[Service]
ExecStart=/home/n1ete/bin/backup-packages
